# packaging/gateway/scripts/dev/sdme/etc/postgresql/chan-pg-bootstrap.service
[Unit]
Description=chan-psql dev first-boot bootstrap (role + databases)
After=postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/chan-pg-bootstrap
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
