# gateway/crates/devserver-control/packaging/chan-gateway-devserver-control.service
[Unit]
Description=chan-gateway devserver-control service (proxy fleet control plane)
Documentation=https://github.com/fiorix/chan/tree/main/gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=chan-gateway-devserver-control
Group=chan-gateway-devserver-control
# Shared public-origin config (domain.env) loaded first; the
# per-service file may override. Optional (`-`) so a stack that sets
# the origins inline still starts.
EnvironmentFile=-/etc/chan-gateway/domain.env
EnvironmentFile=/etc/chan-gateway/devserver-control.env
ExecStart=/usr/bin/chan-gateway-devserver-control
Restart=on-failure
RestartSec=2

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
UMask=0077

[Install]
WantedBy=multi-user.target
