# gateway/crates/workspace-proxy/packaging/chan-gateway-workspace-proxy.service
[Unit]
Description=chan-gateway workspace-proxy service (workspace.chan.app)
Documentation=https://github.com/fiorix/chan/tree/main/gateway
After=network-online.target postgresql.service chan-gateway-profile.service
Wants=network-online.target chan-gateway-profile.service

[Service]
Type=simple
User=chan-gateway
Group=chan-gateway
# Shared domain config (CHAN_DOMAIN + PUBLIC_SCHEME) loaded first; the
# per-service file may override. Optional (`-`) so a stack that sets
# the domain inline still starts.
EnvironmentFile=-/etc/chan-gateway/domain.env
EnvironmentFile=/etc/chan-gateway/workspace-proxy.env
ExecStart=/usr/bin/chan-gateway-workspace-proxy
Restart=on-failure
RestartSec=2

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
UMask=0077

[Install]
WantedBy=multi-user.target
